# root/etc/s6-overlay/s6-rc.d/syncthing/run
#!/usr/bin/execlineb -P

with-contenv
exec syncthing --no-browser --home=/var/syncthing --gui-address=0.0.0.0:8384
